# files/monitor.service
[Unit]
Description=My Monitoring Script

[Service]
ExecStart=/home/ubuntu/<NAME-OF-YOUR-MONITORING-SCRIPT>.sh

[Install]
WantedBy=multi-user.target
